# task051/10-add-usb-kbd.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="03", ATTR{bInterfaceProtocol}=="02", RUN+="/sbin/insmod /home/yan/Work/task05/hello.ko"

